# python/crontab
0 1,13 * * * cd /usr/src/app && python /usr/src/app/RunDailyFinance.py >> cron_log.log 2>&1
