# python/crontab
*/1 * * * * /usr/src/run_finance.sh >> /usr/src/app/logs/cron_log.log 2>&1
