# python/crontab
0 1,9,13,19,21 * * * /usr/src/run_finance.sh >> /usr/src/app/logs/cron_log.log 2>&1
